# systemd/xfce-night-switch-startup.service
[Unit]
Description=Apply day/night theme at session start
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/auto-theme.sh

[Install]
WantedBy=default.target
